# .install/systemd/forgeframe-http-helper.service
[Unit]
Description=ForgeFrame ACME helper listener on port 80
After=network-online.target
Wants=network-online.target
ConditionPathExists=/opt/ForgeFrame/.install/etc/forgeframe.env

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/ForgeFrame
Environment=FORGEFRAME_ENV_FILE=/opt/ForgeFrame/.install/etc/forgeframe.env
ExecStart=/opt/ForgeFrame/scripts/start-http-helper.sh
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/opt/ForgeFrame/.install/lib
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
