# deploy/systemd/forgeframe-retention.service
[Unit]
Description=ForgeFrame retention maintenance
After=network-online.target
Wants=network-online.target
ConditionPathExists=@@ENV_FILE@@

[Service]
Type=oneshot
User=@@SYSTEM_USER@@
Group=@@SYSTEM_GROUP@@
WorkingDirectory=@@INSTALL_ROOT@@
Environment=FORGEFRAME_ENV_FILE=@@ENV_FILE@@
ExecStart=@@INSTALL_ROOT@@/.venv/bin/python @@INSTALL_ROOT@@/scripts/run-history-retention.py --apply
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=@@STATE_DIR@@ @@INSTALL_ROOT@@/backend/.forgeframe
StandardOutput=journal
StandardError=journal
